# x86_64/hadoop/hadoop-resourcemanager.service
[Unit]
Description=Haoop resource manager server

[Service]
Type=forking
User=hadoop
Group=hadoop
EnvironmentFile=/etc/conf.d/hadoop
ExecStart=/usr/bin/yarn --daemon start resourcemanager
ExecStop=/usr/bin/yarn --daemon stop resourcemanager
Restart=on-abnormal
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
